# docker-files/crontabs/www-data
# m h dom mon dow  command

# Invoice generator
*/5 * * * * /var/www/bin/console app:invoice:creator -vvv
